# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[program:bot]
command=python3 bot.py
directory=/app     # Changed from '.' to '/app' since Heroku uses /app as root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0    # Added to prevent log rotation issues
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0    # Added to prevent log rotation issues
environment=PYTHONUNBUFFERED=1   # Added for better logging
